# deploy/systemd/netify-fwa.service.in
[Unit]
Description=Netify FWA - Firewall Agent
Before=network-pre.target
Wants=network-pre.target
After=firewalld.service

[Service]
Type=forking
PIDFile=@volatilestatedir@/netify-fwa.pid
User=root
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=@sbindir@/netify-fwa
KillMode=mixed
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
