# deploy/pfsense/files/netify-fwa.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: netify-fwa
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# netify_fwa_enable (bool):	Set to NO by default.
#				Set it to YES to enable netify-fwa.

. /etc/rc.subr

name="netify_fwa"
rcvar="netify_fwa_enable"

procname="/usr/local/bin/python3.6"
command="%%PREFIX%%/sbin/netify-fwa"
pidfile="/var/run/netify-fwa/netify-fwa.pid"

load_rc_config netify_fwa

: ${netify_fwa_enable:=no}

run_rc_command "$1"
